# com.user.surfmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.surfmonitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>/Users/p5ai/surf_monitor/surf_monitor.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/p5ai/surf_monitor</string>

  <!-- Start at login and restart if it crashes -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Throttle restarts to avoid tight crash loops -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/p5ai/.surf_monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/p5ai/.surf_monitor.log</string>
</dict>
</plist>
